# packaging/debian/prerm
#!/bin/bash
# Pre-remove script for kidlock
set -e

SERVICE_NAME="kidlock"

# Stop and disable the service
if systemctl is-active --quiet "$SERVICE_NAME"; then
    systemctl stop "$SERVICE_NAME" || true
fi

if systemctl is-enabled --quiet "$SERVICE_NAME" 2>/dev/null; then
    systemctl disable "$SERVICE_NAME" || true
fi

exit 0
